# deploy/systemd/alpha-pos-autodeploy.service
# Oneshot unit that runs the self-redeploy poll once. Driven by the matching
# .timer. Install: see deploy/DEPLOY.md ("Self-redeploy without CI").
#
# Adjust User= and the ALPHA_DIR/ExecStart path below to your repo location
# (default: /var/www/alpha_pos). Verify with: ls -d <path>/.git
[Unit]
Description=Alpha POS self-redeploy (git pull + docker compose up --build)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Environment=DEPLOY_BRANCH=main
Environment=ALPHA_DIR=/var/www/alpha_pos
ExecStart=/usr/bin/bash /var/www/alpha_pos/deploy/auto_redeploy.sh
# A failed public-route check must remain visible to operators and monitoring.
# The timer will invoke a fresh oneshot later; auto_redeploy.sh retries the
# current revision when Git is already at the remote SHA but the route is bad.
SuccessExitStatus=0
